# test-owdev/services.d/owdev/finish
#!/bin/sh

echo "Stopping the cluster..."
cd /home/ow/
make stop